# ops/ansible/roles/bfd-server/templates/bluebutton-appserver.service.j2
[Unit]
Description=BFD Server
Wants=network-online.target
After=network.target network-online.target

[Service]
WorkingDirectory={{ data_server_dir }}
ExecStart={{ data_server_dir }}/bfd-server.sh
User={{ data_server_user }}

[Install]
WantedBy=multi-user.target
