# ops/ansible/roles/bfd-server/templates/bluebutton-appserver.service.j2
[Unit]
Description=BFD Server
Wants=network-online.target
After=network.target network-online.target

[Service]
WorkingDirectory={{ systemd/working_dir }}
ExecStart={{ systemd/exec_start }}
User={{ systemd/user }}

[Install]
WantedBy=multi-user.target
