# bluebutton-data-pipeline-benchmarks/src/test/ansible/roles/server_etl/templates/bluebutton-data-pipeline-app.service.j2
[Unit]
Description=Blue Button ETL Service

[Service]
WorkingDirectory=/usr/local/bluebutton-data-pipeline-app
ExecStart=/usr/local/bluebutton-data-pipeline-app/bluebutton-data-pipeline-app.sh
User=bb-etl

[Install]
WantedBy=multi-user.target
