# ops/ansible/roles/bfd-db-migrator/templates/bfd-db-migrator.service.j2
[Unit]
Description=BFD DB Migrator Service
Wants=network-online.target

[Service]
WorkingDirectory={{ working_dir }}
ExecStart={{ exec_start }}
User={{ user }}

[Install]
WantedBy=multi-user.target
