# ops/ansible/roles/bfd-db-migrator/templates/bfd-db-migrator.service.j2
[Unit]
Description=BFD DB Migrator Service
Wants=network-online.target
After=network.target network-online.target {% if (migrator_monitor_enabled|default('', false)) and (bfd_env|default('', false)) %}migrator-monitor.service{% endif %}

{% if (migrator_monitor_enabled|default('', true)) and (bfd_env|default('', true)) %}Requires=migrator-monitor.service{% endif %}


[Service]
WorkingDirectory={{ db_migrator_dir }}
ExecStart={{ db_migrator_dir }}/bfd-db-migrator-service.sh
User={{ db_migrator_user }}

[Install]
WantedBy=multi-user.target
