# apps/bfd-pipeline/bfd-pipeline-benchmarks/src/test/ansible/roles/server_fhir/templates/bluebutton-server-app.service.j2
[Unit]
Description=Blue Button FHIR Server

[Service]
WorkingDirectory=/usr/local/bluebutton-server-app
ExecStart=/usr/local/bluebutton-server-app/wildfly-8.1.0.Final/bin/standalone.sh
User=bb-fhir

[Install]
WantedBy=multi-user.target
